# deploy/refprop8002.service
[Unit]
Description=REFPROP 8002 后端（供 refbackend wine 客户端调用）
After=network.target

[Service]
Type=simple
WorkingDirectory=/www/refprop/refbackend
Environment="PATH=/www/refprop/refbackend/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/www/refprop/refbackend/deploy/refbackend.env
ExecStart=/www/refprop/refbackend/venv/bin/python -m uvicorn wine_refprop_backend:app --host 127.0.0.1 --port 8002
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
